# MVS_Camera/Python/GrabImage.service


[Unit]
Description=StartGrab

[Service]
Type=simple
ExecStart=/bin/bash StartGrab.sh


[Install]
WantedBy=multi-user.target
